# chezmoi/dot_config/autostart/kitty.desktop
[Desktop Entry]
Type=Application
Name=kitty
Exec=kitty
Icon=kitty
X-GNOME-Autostart-enabled=true
